# 04 Ressourcen/Scripts/bridge-stability/bridge-stability-monitor.service
[Unit]
Description=Hybrid-Bridge Stability Snapshot (oneshot)

[Service]
Type=oneshot
ExecStart=%h/livediag_bridge_stability.sh
